# roles/xssproxy/templates/xssproxy.service.j2
[Unit]
Description=XSS Proxy (D-Bus inhibitor)

[Service]
ExecStart={{ xssproxy_bin }}
Restart=on-failure

[Install]
WantedBy=default.target
